# vmname.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/vmname.sh

[Install]
WantedBy=default.target